# rootfs/etc/services.d/crafty_web/run
#!/usr/bin/execlineb -P
cd /crafty_web
s6-setuidgid crafty
/usr/local/bin/python /crafty_web/crafty.py -c "/crafty_web/configs/docker_config.yml"
